# frontend/deployments/nginx-supervisord.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
priority=999
stderr_logfile=/var/log/supervisor/nginx-error.log
stdout_logfile=/var/log/supervisor/nginx-access.log
